# example/pelagius.service
[Unit]
Description=Pelagius
After=network.target

[Service]
Type=simple
Restart=always
WorkingDirectory=/path/to/bot
ExecStart=pnpm start
Environment=NODE_VERSION=22
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
